# iii.conf
[program:iii]
command=python3 /home/curator2/iii/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator2/iii/test.err.log
stdout_logfile=/home/curator2/iii/test.out.log